# server/deploy/camera-server.service
[Unit]
Description=Camera cloud backend (gunicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/camera-server
Environment=HOME=/opt/camera-server
EnvironmentFile=/opt/camera-server/.env
ExecStart=/opt/camera-server/venv/bin/gunicorn \
    --workers 2 --threads 4 --worker-class gthread --timeout 120 \
    --bind 127.0.0.1:8090 server:app
# graceful reload: master keeps the listen socket, workers respawn with new code
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
