# server/deploy/camera-retention.timer
[Unit]
Description=Run camera retention cleanup daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
